# files/macos/com.spkane.calblink.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
	"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.spkane.calblink</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/env</string>
		<string>calblink</string>
		<string>-runAsService</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>ProcessType</key>
	<string>Background</string>

	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>

	<key>ThrottleInterval</key>
	<integer>30</integer>

	<key>StandardOutPath</key>
	<string>/tmp/com.spkane.calblink.out.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/com.spkane.calblink.err.log</string>
</dict>
</plist>
